# deploy/supervisord.conf
; WeChatBot v2.0 - Supervisor config
; Place in /etc/supervisor/conf.d/wechat-bot.conf

[program:wechat-bot]
command=%(here)s/.venv/bin/python %(here)s/main.py -c %(here)s/config.yaml
directory=%(here)s
user=wechat-bot
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=10
stopasgroup=true
killsasgroup=true

stdout_logfile=%(here)s/data/supervisor_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

stderr_logfile=%(here)s/data/supervisor_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
